# root/etc/services.d/tail-error/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc tail -F /config/data/logs/error.log
